# deploy/linux/cammonitor.service
[Unit]
Description=CamMonitor Local CCTV Dashboard
After=network-online.target go2rtc.service
Wants=network-online.target

[Service]
Type=simple
User=cammonitor
Group=cammonitor
WorkingDirectory=/opt/cctv-monitor-kumarans
Environment=PYTHONUNBUFFERED=1
Environment=CAMMONITOR_STRICT_STARTUP=1
ExecStart=/opt/cctv-monitor-kumarans/.venv/bin/python /opt/cctv-monitor-kumarans/run_local.py
Restart=always
RestartSec=5
UMask=007

[Install]
WantedBy=multi-user.target
