# deploy/linux/go2rtc.service
[Unit]
Description=go2rtc Stream Bridge for CamMonitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cammonitor
Group=cammonitor
WorkingDirectory=/opt/cctv-monitor-kumarans
ExecStart=/usr/local/bin/go2rtc -config /opt/cctv-monitor-kumarans/deploy/linux/go2rtc.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
