# deploy/linux/go2rtc.service
[Unit]
Description=go2rtc Stream Bridge for CamMonitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cammonitor
Group=cammonitor
WorkingDirectory=/opt/cam-monitor-v2
ExecStart=/usr/local/bin/go2rtc -config /opt/cam-monitor-v2/deploy/linux/go2rtc.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
